# stage3/05-install-oxigotchi/files/etc/systemd/system/bt-agent.service
[Unit]
Description=Bluetooth Auto-Pair Agent
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=dbus
BusName=org.pwnagotchi.BTAgent
ExecStart=/usr/local/bin/oxigotchi-bt-agent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=bluetooth.target