# stage3/05-install-oxigotchi/files/etc/systemd/system/oxigotchi.service
[Unit]
Description=Pwnagotchi Zero Core Daemon
After=network-online.target bluetooth.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
# The daemon does not implement sd_notify, so Type=simple (not notify) and no
# WatchdogSec — otherwise systemd waits for a readiness signal that never
# arrives, kills the process on timeout, and gives up after a few retries.
Type=simple
ExecStart=/usr/local/bin/oxigotchi
Restart=on-failure
RestartSec=5
Nice=-5
LimitNOFILE=65536
MemoryMax=200M
CPUQuota=80%
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target