# stage3/05-install-oxigotchi/files/etc/systemd/system/usb-gadget.service
[Unit]
Description=oxigotchi USB RNDIS ethernet gadget
DefaultDependencies=no
After=sys-kernel-config.mount systemd-modules-load.service
Before=NetworkManager.service network-pre.target
Wants=sys-kernel-config.mount

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/usb-gadget.sh

[Install]
WantedBy=sysinit.target
